# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/sites/SITENAME/source

exec ../venv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    --workers 4 \
    mngweb.wsgi:application
